# annad.service
# Anna Assistant Daemon - Systemd Service Unit
#
# Anna is a friendly, intelligent system administrator for Arch Linux
# that monitors your system and provides personalized recommendations.
#
# Installation:
#   sudo cp annad.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now annad
#
# Management:
#   sudo systemctl status annad    # Check status
#   sudo systemctl restart annad   # Restart daemon
#   sudo systemctl stop annad      # Stop daemon
#   journalctl -u annad -f         # View logs
#
# Configuration:
#   Use 'annactl config' to manage Anna's settings
#
# Documentation:
#   https://github.com/jjgarcianorway/anna-assistant/blob/main/README.md

[Unit]
Description=Anna Assistant - Autonomous System Administrator for Arch Linux
Documentation=https://github.com/jjgarcianorway/anna-assistant
Documentation=https://wiki.archlinux.org/
After=network.target local-fs.target

# Ensure critical system services are available
Wants=systemd-journald.service

# Restart rate limiting (moved from [Service] - rc.13.1 fix)
StartLimitIntervalSec=0

[Service]
Type=simple

# Phase 3.9: Run as dedicated service user (not root)
# The annad user should be created during package installation
# For now, keep root for compatibility but restrict capabilities
User=root
Group=anna
SupplementaryGroups=

# Phase 3.9: Capability restrictions
# Only grant minimum necessary capabilities
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_CHOWN CAP_FOWNER CAP_SYS_ADMIN CAP_SYS_RESOURCE
SecureBits=keep-caps

# rc.13.3: Ensure runtime directory exists with correct ownership before startup
# This fixes race condition where socket creation fails if /run/anna doesn't exist
PermissionsStartOnly=true
ExecStartPre=/usr/bin/install -d -m0750 -o root -g anna /run/anna
ExecStartPre=/bin/rm -f /run/anna/anna.sock
ExecStart=/usr/local/bin/annad

# Restart policy - automatically recover from failures (RC.7)
Restart=on-failure
RestartSec=10s
StartLimitBurst=3

# Phase 0.4: Security hardening - comprehensive protection
# Prevent privilege escalation
NoNewPrivileges=true

# Filesystem protection
ProtectSystem=strict
ProtectHome=yes
# Beta.115: Allow writing to /usr/local/bin for auto-updates
ReadWritePaths=/var/log/anna /var/lib/anna /run/anna /usr/local/bin

# Use private /tmp directory
PrivateTmp=true

# Drop all ambient capabilities
AmbientCapabilities=

# Phase 3.9: Kernel and system protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Phase 3.9: Network restrictions
# NOTE: Removed IPAddressDeny/IPAddressAllow due to MediaTek mt7921e driver bug
# These restrictions cause netlink operations that crash the WiFi driver
# See: "nl80211: send_event_marker failed: Source based routing not supported"
# v10.1.1: Added AF_NETLINK to allow 'ip' command to work (uses netlink sockets)
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
# IPAddressDeny=any
# IPAddressAllow=localhost

# Phase 3.9: System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @obsolete @resources
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Phase 3.9: Device access control
DevicePolicy=closed
DeviceAllow=/dev/null rw
DeviceAllow=/dev/zero rw
DeviceAllow=/dev/urandom r

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=annad
SyslogLevel=info

# Resource limits
# File descriptors for filesystem watching and IPC
LimitNOFILE=65536

# Nice level (0 = default priority)
Nice=0

# Working directory for daemon (rc.13.2: use root since /var/lib/anna may not exist yet)
WorkingDirectory=/

# Phase 3.9.1: Secure directory structure with anna group write access
# Runtime directory for socket (root:anna 0770)
RuntimeDirectory=anna
RuntimeDirectoryMode=0770

# Ensure files created by service default to 0660 for group anna
UMask=007

# Log directory (root:anna 0750) - group readable for diagnostics
LogsDirectory=anna
LogsDirectoryMode=0750

# State directory (root:anna 0770) - group writable for reports
# Fixes: v3.9.1 - Allow anna group to write health/doctor reports
StateDirectory=anna anna/reports anna/alerts
StateDirectoryMode=0770

# Configuration directory
ConfigurationDirectory=anna
ConfigurationDirectoryMode=0700

[Install]
WantedBy=multi-user.target
