# annad.service
# Anna Assistant Daemon - Systemd Service Unit
#
# Anna is a friendly, intelligent system administrator for Arch Linux
# that monitors your system and provides personalized recommendations.
#
# Installation:
#   sudo cp annad.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now annad
#
# Management:
#   sudo systemctl status annad    # Check status
#   sudo systemctl restart annad   # Restart daemon
#   sudo systemctl stop annad      # Stop daemon
#   journalctl -u annad -f         # View logs
#
# Configuration:
#   Use 'annactl config' to manage Anna's settings
#
# Documentation:
#   https://github.com/jjgarcianorway/anna-assistant/blob/main/README.md

[Unit]
Description=Anna Assistant - Autonomous System Administrator for Arch Linux
Documentation=https://github.com/jjgarcianorway/anna-assistant
Documentation=https://wiki.archlinux.org/
After=network.target local-fs.target

# Ensure critical system services are available
Wants=systemd-journald.service

# Restart rate limiting (moved from [Service] - rc.13.1 fix)
StartLimitIntervalSec=0

[Service]
Type=simple

# Run as root (required for system administration)
# but use anna group for socket access control
User=root
Group=anna
SupplementaryGroups=anna

# rc.13.3: Ensure runtime directory exists with correct ownership before startup
# This fixes race condition where socket creation fails if /run/anna doesn't exist
PermissionsStartOnly=true
ExecStartPre=/usr/bin/install -d -m0750 -o root -g anna /run/anna
ExecStartPre=/bin/rm -f /run/anna/anna.sock
ExecStart=/usr/local/bin/annad

# Restart policy - automatically recover from failures (RC.7)
Restart=on-failure
RestartSec=10s
StartLimitBurst=3

# Phase 0.4: Security hardening - comprehensive protection
# Prevent privilege escalation
NoNewPrivileges=true

# Filesystem protection
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/anna /run/anna

# Use private /tmp directory
PrivateTmp=true

# Drop all ambient capabilities
AmbientCapabilities=

# Additional security restrictions
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=annad
SyslogLevel=info

# Resource limits
# File descriptors for filesystem watching and IPC
LimitNOFILE=65536

# Nice level (0 = default priority)
Nice=0

# Working directory for daemon (rc.13.2: use root since /var/lib/anna may not exist yet)
WorkingDirectory=/

# Phase 0.4: Secure directory structure
# Runtime directory for socket (root:anna 0750)
RuntimeDirectory=anna
RuntimeDirectoryMode=0750

# Ensure files created by service default to 0660 for group anna (rc.13.1 fix)
UMask=007

# Log directory (root:root 0700)
LogsDirectory=anna
LogsDirectoryMode=0700

# State directory (root:root 0700)
# rc.13.2: Ensure /var/lib/anna and subdirectories exist
StateDirectory=anna anna/reports anna/alerts
StateDirectoryMode=0700

# Configuration directory
ConfigurationDirectory=anna
ConfigurationDirectoryMode=0700

[Install]
WantedBy=multi-user.target
