# systemd/annad.service
[Unit]
Description=Anna Daemon - Evidence Oracle
Documentation=https://github.com/jjgarcianorway/anna-assistant
After=network.target

[Service]
Type=simple
User=anna
Group=anna
ExecStart=/usr/local/bin/annad
WorkingDirectory=/usr/share/anna
Restart=on-failure
RestartSec=5
TimeoutStartSec=10
TimeoutStopSec=10

# Environment
Environment=RUST_LOG=annad=info

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

# Allow network (localhost only)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Allow reading system info
ReadOnlyPaths=/proc /sys
ReadWritePaths=/var/lib/anna /var/log/anna /run/anna

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
