# systemd/anna-daemon.service
[Unit]
Description=Anna Assistant Daemon
Documentation=https://github.com/yourusername/anna-assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=anna
Group=anna
WorkingDirectory=/var/lib/anna

# Environment
Environment=ANNA_STATE_DIR=/var/lib/anna
Environment=ANNA_LOG_DIR=/var/log/anna
Environment=RUST_LOG=info

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/anna /var/log/anna
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

# Execution
ExecStart=/usr/bin/annad
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s
TimeoutStopSec=30s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anna-daemon

[Install]
WantedBy=multi-user.target
