# packaging/deb/postinst
#!/bin/bash
# Debian postinst script for anna-daemon
# Runs after package installation

set -e

case "$1" in
    configure)
        # Create system user and group
        if ! getent group anna >/dev/null 2>&1; then
            addgroup --system anna
        fi

        if ! getent passwd anna >/dev/null 2>&1; then
            adduser --system \
                --ingroup anna \
                --home /var/lib/anna \
                --no-create-home \
                --shell /usr/sbin/nologin \
                --gecos "Anna Assistant Daemon" \
                anna
        fi

        # Create directories
        mkdir -p /var/lib/anna /var/log/anna
        chmod 750 /var/lib/anna /var/log/anna
        chown anna:anna /var/lib/anna /var/log/anna

        # Create README files
        cat > /var/lib/anna/README <<'EOF'
Anna Assistant State Directory
===============================

This directory contains persistent state for the Anna daemon.
See /usr/share/doc/anna-daemon/README for details.

WARNING: Do not manually edit state files while daemon is running.
EOF

        cat > /var/log/anna/README <<'EOF'
Anna Assistant Log Directory
============================

This directory contains append-only logs for the Anna daemon.
See /usr/share/doc/anna-daemon/README for monitoring commands.
EOF

        chown anna:anna /var/lib/anna/README /var/log/anna/README
        chmod 640 /var/lib/anna/README /var/log/anna/README

        # Reload systemd
        systemctl daemon-reload || true

        # Note: Service is NOT enabled by default
        # Operator must manually enable after reviewing configuration

        echo ""
        echo "Anna Assistant installed successfully."
        echo ""
        echo "The service is NOT enabled by default."
        echo "Review the configuration and then:"
        echo "  systemctl enable anna-daemon"
        echo "  systemctl start anna-daemon"
        echo ""
        echo "See /usr/share/doc/anna-daemon/PRODUCTION_DEPLOYMENT.md for details."
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
