# etc/systemd/annad.service
[Unit]
Description=Anna v0.11.0 Event-Driven Intelligence Daemon
Documentation=https://github.com/anna-assistant/anna
After=network.target local-fs.target

[Service]
Type=simple

# Run as anna user (unprivileged observer)
# Event-driven architecture: observe, analyze, alert (no auto-repair yet)
User=anna
Group=anna

# Capability Bounding Set (read-only observer, no privileged operations)
# Running as anna user - no capabilities needed for observation
AmbientCapabilities=
CapabilityBoundingSet=

# Execution
ExecStart=/usr/local/bin/annad
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource Limits
CPUQuota=5%
MemoryMax=80M
TasksMax=32

# Managed Directories (systemd auto-creates and sets ownership)
StateDirectory=anna
StateDirectoryMode=0750
LogsDirectory=anna
LogsDirectoryMode=0750
RuntimeDirectory=anna
RuntimeDirectoryMode=0750
UMask=007

# Security Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ProtectKernelLogs=yes

# Writable Paths (anna's own state directories only)
ReadWritePaths=/var/lib/anna /var/log/anna /run/anna
# Read-only access to system state for observation
ReadOnlyPaths=/var/lib/pacman /etc /proc/self/mountinfo

# Need access to /dev for hardware sensors
PrivateDevices=no

[Install]
WantedBy=multi-user.target
