# etc/systemd/annad.service
[Unit]
Description=Anna Assistant Daemon
Documentation=https://github.com/anna-assistant/anna
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=anna
Group=anna
ExecStart=/usr/local/bin/annad
RuntimeDirectory=anna
RuntimeDirectoryMode=0750
StateDirectory=anna
StateDirectoryMode=0750
CacheDirectory=anna
CacheDirectoryMode=0750
Restart=on-failure
RestartSec=3s

# Resource accounting and limits
CPUAccounting=true
MemoryAccounting=true
CPUQuota=50%
MemoryMax=300M
WatchdogSec=60s

# Security hardening
PrivateTmp=yes
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=annad

[Install]
WantedBy=multi-user.target
